# felix-debian11/felix/recipes/attic/x-s-fix_pulse_audio_no_sound_after_suspend/fix_pulse_audio_no_sound_after_suspend.service
[Unit]
Description=Fix PulseAudio no sound after suspend
After=suspend.target

[Service]
Type=oneshot
ExecStart=/usr/bin/pasuspender /bin/true

[Install]
WantedBy=suspend.target
