# 17.04_ZestyZapus/system/0220-s-configure_acpi_wakeup/configure_acpi_wakeup.service
[Unit]
Description=Configure ACPI wakeup

[Service]
ExecStart=/bin/bash -c "echo LID0 >> /proc/acpi/wakeup;"

[Install]
WantedBy=multi-user.target
