# felix-debian12/recipes/03-hw_specifics/macbook_air_7_1/xxxx-s-configure_acpi/fix_macbook_air_wakes_up_after_sleep.service
[Unit]
Description=Fix macbook air wakes up after sleep

[Service]
Type=oneshot
ExecStart=/bin/sh -c "echo XHC1 > /proc/acpi/wakeup && echo LID0 > /proc/acpi/wakeup"

[Install]
WantedBy=multi-user.target

