# 18.04_BionicBeaver/system/0220-s-configure_acpi_wakeup/configure_acpi_wakeup.service
[Unit]
Description=Configure ACPI wakeup

[Service]
ExecStart=

[Install]
WantedBy=multi-user.target
