# felix-ubuntu1804/felix/recipes/attic/x-s-fix_drivers_for_realtek_ethernet_controllers/r8168_fix_after_suspend.service
[Unit]
Description=Insert r8169 module after suspend/hybrid-sleep/hibernate
After=suspend.target
After=hybrid-sleep.target
After=hibernate.target

[Service]
Type=simple
ExecStart=/usr/bin/modprobe r8168

[Install]
WantedBy=suspend.target
WantedBy=hybrid-sleep.target
WantedBy=hibernate.target
