# felix-debian12/user-dotfiles/.config/autostart/screensaver-settings.desktop
[Desktop Entry]
Type=Application
Name=Screensaver
Comment=Set screensaver timeouts
Exec=xset s 0 dpms 600 600 600
