# ubuntu_1804/specifics/macbook_air_7_1/0070-s-disable_gpe66/disable_gpe66.service
[Unit]
Description=Disable GPE 66 causing battery draining on Mac

[Service]
ExecStart=/usr/bin/env bash -c 'echo "disable" > /sys/firmware/acpi/interrupts/gpe66'

[Install]
WantedBy=multi-user.target
